# change-wallpaper.service
[Unit]
Description=Change wallpaper based on time of day

[Service]
Type=oneshot
ExecStart=%h/.local/bin/change_wallpaper.sh