# systemd/generate-wallpaper-timer.path
[Unit]
Description=Watch schedule.json and regenerate wallpaper timer on change

[Path]
PathModified=%h/.config/wallpaper-schedule/schedule.json

[Install]
WantedBy=default.target
